# application_boot.service
[Unit]
Description=to start service on reboot

[Service]
Type=simple
ExecStart=/bin/sh /tmp/node_startup.sh

[Install]
WantedBy=multi-user.target
